# deploy/dtm-ai.service.d/hermes-rw.conf
# systemd drop-in: let the dtm-ai web service WRITE Hermes' data dir, so the owner can swap
# Hermes' brain (cloud↔local) and later edit profiles/SOUL/memory from the dashboard.
# ProtectSystem=strict makes the FS read-only except ReadWritePaths; this ADDS /srv/hermes-data
# to the existing /opt/dtm-ai (ReadWritePaths is additive across drop-ins) without touching the
# installed unit (which runs --port 8090).
#
# Install:
#   sudo mkdir -p /etc/systemd/system/dtm-ai.service.d
#   sudo cp /opt/dtm-ai/deploy/dtm-ai.service.d/hermes-rw.conf /etc/systemd/system/dtm-ai.service.d/
#   sudo systemctl daemon-reload && sudo systemctl restart dtm-ai
[Service]
ReadWritePaths=/srv/hermes-data
