# deploy/msp-ai.service.d/10-full-access.conf
# D-22 — drop-in that relaxes the msp-ai systemd sandbox so the admin Terminal's `sudo` actually works.
#
# The base unit sets NoNewPrivileges=true, RestrictSUIDSGID=true, ProtectSystem=strict (+ ProtectHome,
# ReadWritePaths). Those three each BLOCK sudo: no privilege escalation, no setuid binaries, read-only FS.
# This drop-in turns them off so a logged-in admin can make real system changes from the :8090 Terminal.
#
# Install (run as root, ONE line at a time):
#   sudo mkdir -p /etc/systemd/system/msp-ai.service.d
#   sudo cp /opt/msp-ai/deploy/msp-ai.service.d/10-full-access.conf /etc/systemd/system/msp-ai.service.d/
#   sudo systemctl daemon-reload
#   sudo systemctl restart msp-ai
#
# ⚠ TRADE-OFF: this removes most of the always-on web app's OS-level isolation. If you'd rather keep the
#   main app hardened, SKIP this file and use the :8091 recovery console (runs as root) for system changes.
#   Remove with: sudo rm /etc/systemd/system/msp-ai.service.d/10-full-access.conf && sudo systemctl daemon-reload && sudo systemctl restart msp-ai
[Service]
NoNewPrivileges=false
RestrictSUIDSGID=false
ProtectSystem=no
ProtectHome=no
PrivateTmp=false
ProtectKernelTunables=false
ProtectControlGroups=false
LockPersonality=false
ReadWritePaths=
